# launchd/com.frontier.dealcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.frontier.dealcheck</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/muham/Projects/Frontier-GoWild-Search/gowild_deal_report.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/muham/Projects/Frontier-GoWild-Search</string>

    <!-- Tue(2), Wed(3), Thu(4) at 00:01 local time -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>0</integer>
            <key>Minute</key><integer>1</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>0</integer>
            <key>Minute</key><integer>1</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>0</integer>
            <key>Minute</key><integer>1</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/muham/Projects/Frontier-GoWild-Search/results/dealcheck.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/muham/Projects/Frontier-GoWild-Search/results/dealcheck.err.log</string>
</dict>
</plist>
